# config/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=nobody:nogroup
 
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
 
[inet_http_server]
port = 9001
username = root
password = root
 
#[program:php]
#command=docker-php-entrypoint php-fpm # This is our main process so let us put it here
 
[include]
files=/etc/supervisor/conf.d/*.conf

[program:vk]
command=php /var/www/lenta.events/app/jobs/vk.php -DFOREGROUND
directory=/opt
autostart=true
autorestart=true
startretries=5
startsecs=0
user=root
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB